# zero2/zero2/zero2-mqtt.service
[Unit]
Description=Zero2 Sensor MQTT Publisher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sugi
WorkingDirectory=/home/sugi
ExecStart=/home/sugi/z2env/bin/python /home/sugi/zero2_sensors_mqtt.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
